# templates/default/resque-web-upstart.conf.erb
description "resque web server"

start on runlevel [2345]
stop on runlevel [016]

script
  exec sudo -u <%= node[:rubber_resque][:app_user] %> /bin/bash -l -c 'cd <%= node[:rubber_resque][:app_path] %> && exec bundle exec puma --port <%= node[:rubber_resque][:resque_web][:port] %> config/resque-web.ru &>> log/resque-web.log'
end script

post-start script
  status resque-web | head -n1 | awk '{print $NF}' > <%= node[:rubber_resque][:resque_web][:pid_file] %>
end script

post-stop script
  rm -f <%= node[:rubber_resque][:resque_web][:pid_file] %>
end script
